# etc/profile.d/nc_server.csh
if ( "${path}" !~ */opt/nc_server/bin* ) then
   set path = ( $path /opt/nc_server/bin )
endif
