# etc/profile.d/nc_server.sh
if ! echo ${PATH} | /bin/grep -q /opt/nc_server/bin ; then
   PATH=${PATH}:/opt/nc_server/bin
fi
